# server/observatorio.timer
[Unit]
Description=Dispara observatorio-publish cada 10s

[Timer]
OnBootSec=15
OnUnitActiveSec=10
AccuracySec=2s

[Install]
WantedBy=timers.target
