# server/observatorio.service
[Unit]
Description=Publica observatorio_state (nowplaying + system) en Supabase
After=network-online.target

[Service]
Type=oneshot
# La service key vive aquí (chmod 600), nunca en el repo ni en la web:
EnvironmentFile=%h/.config/observatorio/env
ExecStart=%h/.local/bin/observatorio-publish.sh
Nice=10
