# server/observatorio.timer
[Unit]
Description=Dispara observatorio-publish cada 30s

[Timer]
OnBootSec=20
OnUnitActiveSec=30
AccuracySec=5s

[Install]
WantedBy=timers.target
